# v2ray/services.conf
[program:v2ray]
command=/usr/bin/v2ray/v2ray -config /etc/v2ray/server-config.json
