# shadowsocksr-client/services.conf
[program:shadowsocksr]
priority = 100
command = bash -c 'python /opt/shadowsocksr/shadowsocks/local.py -c /opt/shadowsocksr/user-config.json'
