# shadowsocksr-server/services.conf
[program:shadowsocksr]
priority = 100
command = bash -c 'python /opt/shadowsocks/shadowsocks/server.py'
