# shadowsocksr-server/services.conf
[program:shadowsocksr]
priority = 100
command = bash -c 'python /opt/shadowsocksr/shadowsocks/server.py -c /opt/shadowsocksr/user-config.json'
